# system-monitor.service
[Unit]
Description=System Monitor Widget
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/system_monitor_widget.py
Restart=on-failure

[Install]
WantedBy=default.target
